# scripts_debian/custom_wheezy_whonix-gateway/extra-whonix-files/etc/udev/rules.d/99-qubes-whonix.rules
SUBSYSTEMS=="xen", KERNEL=="eth*", ACTION=="add", RUN+="/usr/lib/whonix/setup-ip"
SUBSYSTEMS=="xen", KERNEL=="eth*", ACTION=="online", RUN+="/usr/lib/whonix/setup-ip"
